# docker_deps/services.d/breedbot/run
#!/usr/bin/with-contenv sh

cd /app
pipenv run python /app/bot.py